# build/files/usr/lib/systemd/system/lateralus-auto-upgrade.timer
[Unit]
Description=Daily check for bootc image upgrades

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
