# build/files/usr/lib/systemd/system/lateralus-flatpak-setup.service
[Unit]
Description=Install Lateralus Flatpak applications
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/lateralus-flatpak-setup
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
