# build/files/usr/lib/systemd/system/lateralus-brew-setup.service
[Unit]
Description=Install Homebrew and packages after rebase
After=network-online.target systemd-tmpfiles-setup.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/lateralus-brew-setup
RemainAfterExit=yes
# Network install of ~50 brew packages can take 10-30 min on first boot
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
