# build/files/usr/lib/systemd/system/lateralus-grub-setup.service
[Unit]
Description=Install Lateralus GRUB theme
ConditionPathExists=!/boot/grub2/themes/lateralus/theme.txt
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/lateralus-grub-setup
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
