# build/files/usr/lib/systemd/user/lateralus-user-setup.service
[Unit]
Description=Populate user dotfiles from /etc/skel after rebase

[Service]
Type=oneshot
ExecStart=/usr/libexec/lateralus-user-setup
RemainAfterExit=yes

[Install]
WantedBy=default.target
